# Monitor/svcctl-monitor.service
[Unit]
Description=Service Monitor (systemd watchdog)
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/service-monitor
ExecStart=/usr/bin/python3 /opt/service-monitor/monitor.py --config /opt/service-monitor/config.json
Restart=on-failure
RestartSec=5
User=svcctl
Group=svcctl
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
